# taut-server.service
[Unit]
Description=AKA Python Youtube-taut tornado server service

[Service]
Type=simple
WorkingDirectory=/var/www/youtube-taut/server/
ExecStart=/usr/bin/python server.py --logging=debug
KillMode=process
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target